# deploy/method.service
[Unit]
Description=Method - Research Planner
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/method
EnvironmentFile=/home/ubuntu/method/.env
ExecStart=/home/ubuntu/method/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
